# .config/systemd/user/gdrive-bisync.service
[Unit]
Description=rclone bisync for Google Drive Documents
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/rclone bisync %h/Documents gdrive: \
    --exclude ".git/**" \
    --exclude ".stfolder/**" \
    --exclude "EVE/logs/Gamelogs/**" \
    --drive-import-formats docx,xlsx,pptx,odt,ods \
    --drive-export-formats docx,xlsx,pptx,odt,ods \
    --track-renames \
    --verbose \
    --log-file=%h/.local/state/rclone-bisync.log

