# .config/systemd/user/gdrive-bisync.service
[Unit]
Description=rclone bisync for Google Drive Documents
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/gdrive-bisync-wrapper.sh

